# daemon/pkg/digitalgulag-daemon.service
[Unit]
Description=digitalgulag Activity Tracking Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/digitalgulag-daemon run --headless
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
